# com.sealsuite.vpnmonitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.sealsuite.vpnmonitor</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-l</string>
		<string>/Users/kan.lu/Library/Scripts/sealsuite-vpn-monitor/check_and_reconnect_vpn.sh</string>
	</array>

	<key>RunAtLoad</key>
	<true/>

	<key>StartInterval</key>
	<integer>60</integer>

	<key>StandardOutPath</key>
	<string>/Users/kan.lu/Library/Logs/sealsuite-vpn-monitor-stdout.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/kan.lu/Library/Logs/sealsuite-vpn-monitor-stderr.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
	</dict>

	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>

	<key>ThrottleInterval</key>
	<integer>30</integer>

	<key>ProcessType</key>
	<string>Interactive</string>
</dict>
</plist>
